# kodi-monitor.service
[Unit]
Description=Kodi Monitor Script
After=network.target

[Service]
ExecStart=/usr/bin/python3 /home/USERNAME/kodi-monitor/monitor-kodi.py --log-level INFO --user USERNAME
WorkingDirectory=/home/USERNAME/kodi-monitor/
Restart=always

[Install]
WantedBy=default.target
